# security-stack.service
[Unit]
Description=Security Stack (Pi-hole, Unbound, Squid, Dante)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%h/security-stack

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/security-stack

# Start the stack
ExecStart=/usr/bin/podman-compose up -d

# Stop the stack
ExecStop=/usr/bin/podman-compose down

# Reload
ExecReload=/usr/bin/podman-compose restart

[Install]
WantedBy=default.target
